# donalert.service
[Unit]
Description=DonAlert Donation Platform
After=network.target

[Service]
Type=forking
User=xocoo
Group=xocoo
WorkingDirectory=/srv/www/donalert.invictamotus.com
Environment=PATH=/srv/www/donalert.invictamotus.com/venv/bin
Environment=FLASK_ENV=production
ExecStart=/srv/www/donalert.invictamotus.com/venv/bin/gunicorn --config gunicorn.conf.py wsgi:application --daemon
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
PIDFile=/srv/www/donalert.invictamotus.com/gunicorn.pid
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target