# scripts/donalert-scheduled-changes.service
# DonAlert Scheduled Changes Systemd Service
#
# This service file can be used as an alternative to cron jobs for processing
# scheduled subscription changes.
#
# To install:
#   sudo cp scripts/donalert-scheduled-changes.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable donalert-scheduled-changes.service
#   sudo systemctl start donalert-scheduled-changes.service
#
# To check status:
#   sudo systemctl status donalert-scheduled-changes.service
#
# To view logs:
#   sudo journalctl -u donalert-scheduled-changes.service

[Unit]
Description=DonAlert Scheduled Subscription Changes Processor
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/srv/www/donalert.invictamotus.com
Environment=PATH=/srv/www/donalert.invictamotus.com/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/srv/www/donalert.invictamotus.com/venv/bin/python /srv/www/donalert.invictamotus.com/scripts/process_scheduled_changes.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target